# workflows/aom-buchhaltungsbot@.service
# ──────────────────────────────────────────────────────────────
# aom-buchhaltungsbot@.service
# systemd System-Unit-Template für den Buchhaltungsbot
# (Instanz = Benutzername, z.B. aom-buchhaltungsbot@alice.service)
#
# Installation (als root):
#   cp workflows/aom-buchhaltungsbot@.service /etc/systemd/system/
#   systemctl daemon-reload
#   sudo systemctl enable --now aom-buchhaltungsbot@$(whoami)
#
# Logs anzeigen:
#   journalctl -u aom-buchhaltungsbot@<user> -f
#
# Hintergrund: User-Units können nicht von docker.service abhängen;
# für Headless-Server wird daher eine System-Unit benötigt.
#
# HRFR-v2 · sys-Domain · Stand 2025
# ──────────────────────────────────────────────────────────────

[Unit]
Description=aom Buchhaltungsbot — Rechnungs-RAG + Cloudflare-Tunnel-Verwaltung (%i)
Documentation=https://github.com/aom1941/aom1941/tree/main/workflows
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
# Pfad zu buchhaltungsbot.sh anpassen (Repo-Klonpfad):
ExecStart=/bin/bash /home/%i/aom1941/workflows/buchhaltungsbot.sh run
EnvironmentFile=/home/%i/aom-sys-scanner-rag/.env
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aom-buchhaltungsbot

[Install]
WantedBy=multi-user.target
